# debian/postrm
#! /bin/sh
# postrm script for linuxmuster-ipcop
#
# see: dh_installdeb(1)

set -e

case "$1" in
    remove)
	update-rc.d -f linuxmuster-ipcop remove
	;;

    purge)
	[ -f /etc/init.d/linuxmuster-ipcop ] && rm -f /etc/init.d/linuxmuster-ipcop
	[ -f /etc/default/linuxmuster-ipcop ] && rm -f /etc/default/linuxmuster-ipcop
	[ -d /usr/share/linuxmuster-ipcop ] && rm -rf /usr/share/linuxmuster-ipcop
	[ -d /var/cache/linuxmuster-ipcop ] && rm -rf /var/cache/linuxmuster-ipcop
	[ -d /var/lib/uml/ipcop ] && rm -rf /var/lib/uml/ipcop
	;;

    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
